# kiwix-server/rootfs/etc/services.d/kiwix/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Kiwix
# Shutdown script
# ==============================================================================

bashio::log.info "Stopping Kiwix services..."

# Stop nginx
if [[ -f /tmp/nginx.pid ]]; then
    NGINX_PID=$(cat /tmp/nginx.pid)
    if kill -0 ${NGINX_PID} 2>/dev/null; then
        bashio::log.info "Stopping nginx (PID: ${NGINX_PID})..."
        kill ${NGINX_PID} 2>/dev/null || true
    fi
    rm -f /tmp/nginx.pid
fi

# Stop Kiwix server if running
if [[ -f /tmp/kiwix-serve.pid ]]; then
    KIWIX_PID=$(cat /tmp/kiwix-serve.pid)
    if kill -0 ${KIWIX_PID} 2>/dev/null; then
        bashio::log.info "Stopping Kiwix server (PID: ${KIWIX_PID})..."
        kill ${KIWIX_PID} 2>/dev/null || true
    fi
    rm -f /tmp/kiwix-serve.pid
fi

# Stop management API if running
if [[ -f /tmp/kiwix-manager.pid ]]; then
    MANAGEMENT_PID=$(cat /tmp/kiwix-manager.pid)
    if kill -0 ${MANAGEMENT_PID} 2>/dev/null; then
        bashio::log.info "Stopping management API (PID: ${MANAGEMENT_PID})..."
        kill ${MANAGEMENT_PID} 2>/dev/null || true
    fi
    rm -f /tmp/kiwix-manager.pid
fi
