# kiwix-server/rootfs/etc/services.d/kiwix/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Kiwix
# Runs the Kiwix server and management API behind nginx reverse proxy
# ==============================================================================

bashio::log.info "Starting Kiwix server with nginx reverse proxy..."

# Get configuration values
PORT=$(bashio::config 'port')
ZIM_STORAGE_PATH=$(bashio::config 'zim_storage_path')
ENABLE_MANAGEMENT=$(bashio::config 'enable_management')
LOG_LEVEL=$(bashio::config 'log_level')
MAX_UPLOAD_SIZE=$(bashio::config 'max_upload_size')

# Internal ports (not exposed externally)
KIWIX_INTERNAL_PORT=8080
MANAGEMENT_INTERNAL_PORT=8081

bashio::log.info "Configuration:"
bashio::log.info "  External Port (nginx): ${PORT}"
bashio::log.info "  Kiwix Internal Port: ${KIWIX_INTERNAL_PORT}"
bashio::log.info "  Management Internal Port: ${MANAGEMENT_INTERNAL_PORT}"
bashio::log.info "  ZIM Storage Path: ${ZIM_STORAGE_PATH}"
bashio::log.info "  Management Enabled: ${ENABLE_MANAGEMENT}"
bashio::log.info "  Log Level: ${LOG_LEVEL}"

# Ensure ZIM storage directory exists and has correct permissions
mkdir -p "${ZIM_STORAGE_PATH}"
chown -R kiwix:kiwix "${ZIM_STORAGE_PATH}"
chmod -R g+w "${ZIM_STORAGE_PATH}"

# Start management API in background if enabled
if bashio::var.true "${ENABLE_MANAGEMENT}"; then
    bashio::log.info "Starting management API on internal port ${MANAGEMENT_INTERNAL_PORT}..."
    export ZIM_STORAGE_PATH="${ZIM_STORAGE_PATH}"
    export MAX_UPLOAD_SIZE=${MAX_UPLOAD_SIZE}
    # Start management API - binds to localhost, accessed via nginx proxy
    nohup /opt/venv/bin/python /usr/local/bin/kiwix-manager.py \
        --port ${MANAGEMENT_INTERNAL_PORT} \
        --host 127.0.0.1 \
        --storage-path "${ZIM_STORAGE_PATH}" \
        --max-upload-size ${MAX_UPLOAD_SIZE} \
        > /proc/1/fd/1 2>/proc/1/fd/2 &
    MANAGEMENT_PID=$!
    bashio::log.info "Management API started with PID ${MANAGEMENT_PID}"
    
    # Store PID for potential cleanup
    echo ${MANAGEMENT_PID} > /tmp/kiwix-manager.pid
    
    # Wait a moment for management API to start
    sleep 3
    
    # Check if management API is still running
    if ! kill -0 ${MANAGEMENT_PID} 2>/dev/null; then
        bashio::log.error "Management API failed to start!"
    else
        # Verify it's responding
        if curl -f -s http://localhost:${MANAGEMENT_INTERNAL_PORT}/api/zim > /dev/null 2>&1; then
            bashio::log.info "Management API is responding on port ${MANAGEMENT_INTERNAL_PORT}"
        else
            bashio::log.warning "Management API started but not responding yet. It may need more time."
        fi
    fi
fi

# Check if there are any ZIM files or create empty library
ZIM_COUNT=$(find "${ZIM_STORAGE_PATH}" -maxdepth 1 -name "*.zim" -type f 2>/dev/null | wc -l)
LIBRARY_XML="${ZIM_STORAGE_PATH}/library.xml"

if (( ZIM_COUNT == 0 )); then
    bashio::log.info "No ZIM files found. Creating empty library..."
    # Create an empty library.xml file for kiwix-serve
    cat > "${LIBRARY_XML}" << 'EOF'
<?xml version="1.0" encoding="UTF-8"?>
<library version="2.0" />
EOF
    chown kiwix:kiwix "${LIBRARY_XML}"
    bashio::log.info "Empty library created. Add ZIM files via the management interface."
else
    bashio::log.info "Found ${ZIM_COUNT} ZIM file(s) in storage directory"
fi

# Change to ZIM storage directory
cd "${ZIM_STORAGE_PATH}" || bashio::exit.nok "Could not change to ZIM storage directory"

# Start Kiwix server in background
bashio::log.info "Starting Kiwix server on internal port ${KIWIX_INTERNAL_PORT}..."
bashio::log.info "Kiwix will serve ZIM files from: ${ZIM_STORAGE_PATH}"
bashio::log.info "Using --library and --monitorLibrary for automatic ZIM file detection"

# Note: We can't use --urlRootLocation with a dynamic ingress token since it changes per session
# Instead, we rely on nginx sub_filter and JavaScript to rewrite paths

if [[ -f "${LIBRARY_XML}" ]]; then
    nohup kiwix-serve --library --monitorLibrary --address 127.0.0.1 --port ${KIWIX_INTERNAL_PORT} "${LIBRARY_XML}" \
        > /proc/1/fd/1 2>/proc/1/fd/2 &
else
    nohup kiwix-serve --library --monitorLibrary --address 127.0.0.1 --port ${KIWIX_INTERNAL_PORT} . \
        > /proc/1/fd/1 2>/proc/1/fd/2 &
fi
KIWIX_PID=$!
echo ${KIWIX_PID} > /tmp/kiwix-serve.pid
bashio::log.info "Kiwix server started with PID ${KIWIX_PID}"

# Wait a moment for Kiwix to start
sleep 3

# Verify Kiwix is running
if ! kill -0 ${KIWIX_PID} 2>/dev/null; then
    bashio::log.error "Kiwix server failed to start!"
else
    if curl -f -s http://localhost:${KIWIX_INTERNAL_PORT} > /dev/null 2>&1; then
        bashio::log.info "Kiwix server is responding on port ${KIWIX_INTERNAL_PORT}"
    else
        bashio::log.warning "Kiwix server started but not responding yet. It may need more time."
    fi
fi

# Get IP address for logging (BusyBox compatible)
HA_IP=$(ip route get 1.1.1.1 2>/dev/null | awk '/src/ {print $7}' || \
        ip addr show | grep 'inet ' | grep -v '127.0.0.1' | head -1 | awk '{print $2}' | cut -d/ -f1 || \
        echo "homeassistant-ip")

# Start nginx (foreground - this becomes the main process)
bashio::log.info "Starting nginx reverse proxy on port ${PORT}..."
bashio::log.info "Access URLs:"
bashio::log.info "  Landing page: http://${HA_IP}:${PORT}/"
bashio::log.info "  Kiwix Wiki: http://${HA_IP}:${PORT}/wiki/"
bashio::log.info "  Management: http://${HA_IP}:${PORT}/manage/"
bashio::log.info "  Via Ingress: Access via Home Assistant sidebar (with tabbed interface)"

exec nginx -g "daemon off;"
